# vnc.service
[Service]
User=debian
Restart=always
ExecStart=/home/debian/homelab/vnc.sh
[Install]
WantedBy=multi-user.target
